# systemd/templates/localserver-stack.service.in
[Unit]
Description=__STACK__ compose stack
Documentation=https://github.com/darraghog/localserver-config
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__REPO_ROOT__/scripts/start-stack.sh __STACK__ up
ExecStop=__REPO_ROOT__/scripts/start-stack.sh __STACK__ down

[Install]
WantedBy=default.target
