# systemd/user/localserver-n8n.service
[Unit]
Description=n8n compose stack
Documentation=https://github.com/darraghog/localserver-config
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__REPO_ROOT__/scripts/start-n8n.sh up
ExecStop=__REPO_ROOT__/scripts/start-n8n.sh down

[Install]
WantedBy=default.target
