# systemd/user/localserver-hello-world.service
[Unit]
Description=Hello World compose stack
Documentation=https://github.com/darraghog/localserver-config
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__REPO_ROOT__/compose/hello-world
Environment=PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__HOME__/.local/bin/podman-compose up -d
ExecStop=__HOME__/.local/bin/podman-compose down

[Install]
WantedBy=default.target
